# systemd/bjf-tracker.timer
[Unit]
Description=BJF Tracker (hourly)

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
